# etc/speakeasy.conf
; Sample supervisord configuration for kaldi-gstreamer-server.
; Modify according to you needs and put under /etc/supervisor/conf.d/

[program:speakeasy-full-duplex-english-server]
command = python /usr/lib/kaldi-gstreamer-server/kaldigstserver/master_server.py --port=8000
; user = root
redirect_stderr = true
stdout_logfile = /var/log/speakeasy/server.log
environment= GST_PLUGIN_PATH="/usr/lib/kaldi/src/gst-plugin"
directory = /usr/lib/kaldi-gstreamer-server

[program:speakeasy-full-duplex-english-worker]
command = python /usr/lib/kaldi-gstreamer-server/kaldigstserver/worker.py -u ws://localhost:8000/worker/ws/speech -c sample_worker.yaml
numprocs = 5
process_name=%(program_name)s-%(process_num)s
; user = root
redirect_stderr = true
stdout_logfile = /var/log/speakeasy/worker-%(process_num)s.log
environment= GST_PLUGIN_PATH="/usr/lib/kaldi/src/gst-plugin"
directory = /usr/lib/kaldi-gstreamer-server

